# openskyd-opkg/CONTROL/postinst
echo "Starting openskyd service"

systemctl daemon-reload
systemctl enable openskyd.service
systemctl start openskyd.service

