# openskyd-arch/openskyd.install
post_install() {
	systemctl daemon-reload
	systemctl enable openskyd
}

post_upgrade() {
	systemctl daemon-reload
	systemctl restart openskyd
}

pre_remove() {
	systemctl stop openskyd
	systemctl disable openskyd
}



