# openskyd-arch/openskyd-wd.install
post_install() {
	systemctl daemon-reload
	systemctl enable openskyd-wd
}

post_upgrade() {
	systemctl daemon-reload
	systemctl restart openskyd-wd
}

pre_remove() {
	systemctl stop openskyd
	systemctl disable openskyd-wd
}

