# openskyd-arch/openskyd.install
post_install() {
	echo Enter the serial number: 
	read ser
	echo "Serial = $ser" >> /etc/openskyd.conf
	systemctl daemon-reload
	systemctl enable openskyd
}

